# systemd/ranwall.timer
[Unit]
Description=Periodically randomly change the wallpaper

[Timer]
OnActiveSec=0min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target